# app.conf
[program:app]
command=/bin/bash -lc "cd /opt/brimir && RAILS_ENV=production bundle exec rails s -b 0.0.0.0"
