# docker/app.conf
[program:app]
command=/bin/bash -lc "bash /opt/brimir/boot.bash"
